# attic/20131213/rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user start-up script.
#

hdparm -B 255 /dev/sda
hdparm -B 255 /dev/sdb
