# attic/20161115/profile.d/gl.sh
__GL_SYNC_TO_VBLANK=0
__GL_SYNC_DISPLAY_DEVICE=DFP-0
__VDPAU_NVIDIA_SYNC_DISPLAY_DEVICE=DFP-0

__GLVND_DISALLOW_PATCHING=1

export __GL_SYNC_TO_VBLANK __GL_SYNC_DISPLAY_DEVICE __VDPAU_NVIDIA_SYNC_DISPLAY_DEVICE __GLVND_DISALLOW_PATCHING
