# attic/20180524/systemd/pacman-optimize.service
[Unit]
Description=Optimize pacman database

[Service]
ExecStart=/bin/pacman-optimize

[Install]
WantedBy=basic.target
