# attic/20151231/etc/profile.d/path.sh
PATH=/home/fred/.cabal/bin:$PATH
PATH=/home/fred/Workspace/bin:$PATH
PATH=/home/fred/Workspace/archive/Scripts/System:$PATH

export PATH
