# attic/20151231/home/autostart/xset-dpms.desktop
[Desktop Entry]
Type=Application
Name=xset-dpms
Comment=Load DPMS X Settings
Exec=/usr/bin/xset -dpms
Hidden=false
X-GNOME-Autostart-enabled=true
