# autostart/pcmanfm.desktop
[Desktop Entry]
Name=Desktop Manager
Exec=pcmanfm --desktop
Icon=image-x-generic
Terminal=false
Type=Application
Categories=Utility;
NoDisplay=true
