# attic/20150908/config/autostart/xsetdpms.desktop
[Desktop Entry]
Type=Application
Name=xsetdpms
Comment=Disable DPMS on X
Exec=/usr/bin/xset -dpms
Hidden=false
X-GNOME-Autostart-enabled=true
