# attic/20111212/rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#

hdparm -S 255 -B 255 /dev/sda &>/tmp/hdparm-sda.log
hdparm -S 255 -B 255 /dev/sdb &>/tmp/hdparm-sdb.log

read KEY
