# attic/20191108/udev/wake-on-lan.rules
ACTION=="add", SUBSYSTEM=="net", NAME=="wl*", RUN+="/usr/bin/ethtool -s $name wol g"
