# attic/20180421/systemd/localepurge.service
[Unit]
Description=Purge unneeded locales

[Service]
ExecStart=/bin/localepurge

[Install]
WantedBy=basic.target
