# attic/20200204-rpi/systemd/fileperms.timer
[Unit]
Description=Fix file permissions every hour

[Timer]
OnBootSec=10m
OnUnitActiveSec=5h

[Install]
WantedBy=timers.target
