# attic/20200729/autostart/redshift.desktop
[Desktop Entry]
Version=1.0
Name=Redshift
Exec=redshift
Icon=redshift
Terminal=false
Type=Application
Categories=Utility;
NoDisplay=true
