# attic/20170307-void/rc.local
# Default rc.local for void; add your custom commands here.
#
# This is run by runit in stage 2 before the services are executed
# (see /etc/runit/2).

/usr/bin/hdparm -B 255 -S 255 -W 1 /dev/sda /dev/sdb

